# src/librootfs/akaimer39lib/etc/init.d/rc.local
#!/bin/sh

mount /dev/mmcblk0 /mnt
mount /mnt/local /usr/local

echo ds1307 0x68 > /sys/class/i2c-adapter/i2c-0/new_device
hwclock -s
ifconfig lo up

/usr/bin/tcpsvd 0 21 ftpd -w /mnt &

/usr/bin/telnetd

mkdir -p /etc/dropbear
if [ ! -f /etc/dropbear/dropbear_rsa_host_key ]; then
    dropbearkey -t rsa -f /etc/dropbear/dropbear_rsa_host_key
fi
dropbear -B

ntpd -q -p time.windows.com

#/etc/init.d/wifi
/etc/init.d/monitor.sh &
/etc/jffs2/service.sh &
