# src/librootfs/akaimer39lib/etc/init.d/rc.local
#!/bin/sh

# mount jffs2 file-system. 
/bin/mount -t jffs2 /dev/mtdblock2 /etc
# note: can't recommend running other app before `mount` command.

echo ds1307 0x68 > /sys/class/i2c-adapter/i2c-0/new_device
hwclock -s
ifconfig lo up

/usr/bin/tcpsvd 0 21 ftpd -w /mnt &

/usr/bin/telnetd
ntpd -p 218.75.4.130
/etc/init.d/monitor.sh &
/etc/jffs2/service.sh &

